# root/etc/init.d/wifi-presence
#!/bin/sh /etc/rc.common
# WiFi Presence Monitoring — procd init script
# Manages cron entry for periodic WiFi client checks

START=95
USE_PROCD=1

CRONTAB="/etc/crontabs/root"
CRON_TAG="# wifi-presence"

start_service() {
	local enabled interval

	config_load wifi-presence
	config_get_bool enabled global enabled 0

	# Remove existing entry first
	sed -i "/${CRON_TAG}/d" "$CRONTAB" 2>/dev/null

	if [ "$enabled" -ne 1 ]; then
		/etc/init.d/cron restart
		return
	fi

	config_get interval global interval 1

	echo "*/$interval * * * * /usr/libexec/wifi-presence.sh ${CRON_TAG}" >> "$CRONTAB"
	/etc/init.d/cron restart
}

stop_service() {
	sed -i "/${CRON_TAG}/d" "$CRONTAB" 2>/dev/null
	/etc/init.d/cron restart
}

service_triggers() {
	procd_add_reload_trigger "wifi-presence"
}
